# systemd/instantpostinstall.service
[Unit]
Description=instantOS post installation steps
Before=lightdm.service lightdm-plymouth.service

[Service]
Type=oneshot
ExecStart=/usr/bin/instantpostinstall
KillMode=process
TimeoutStopSec=1000000

[Install]
WantedBy=multi-user.target lightdm.service lightdm-plymouth.service
RequiredBy=lightdm.service lightdm-plymouth.service


