# iso/releng/airootfs/etc/xdg/autostart/instantos-install.desktop
[Desktop Entry]
Type=Application
Name=instantOS Installer
Exec=/usr/local/bin/installapplet
X-GNOME-Autostart-enabled=true
Terminal=false
