# iso/releng/airootfs/etc/xdg/autostart/instantos-liveautostart.desktop
[Desktop Entry]
Type=Application
Name=instantOS Live Tweaks
Exec=/usr/local/bin/liveautostart
X-GNOME-Autostart-enabled=true
NoDisplay=true
Terminal=false
